# rootfs/etc/services.d/cron/run
#!/command/with-contenv bash

printf 'time="%s" level=info msg="%s" type=cron\n' "$(date -u +"%Y-%m-%dT%H:%M:%SZ")" "starting cron"
exec cron -f -L 15
